# httpd.conf/cron.d
* * * * * /etc/rc.d/a/update >/dev/null 2>&1
